# aisst.service
[Unit]
Description=AISST Assistant Service
After=network.target redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/aisst
Environment=PATH="/root/aisst/.venv/bin"
ExecStart=/root/aisst/.venv/bin/gunicorn main:app --workers 2 --worker-class uvicorn.workers.UvicornWorker --bind 0.0.0.0:8000 --timeout 300 --keep-alive 60
ExecStartPre=/root/aisst/.venv/bin/python -m redis_utils.redis_listener &
ExecStartPre=/root/aisst/.venv/bin/python -m rag_chain.rag_worker &
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
